# scripts/my_supervisor.conf
[program: todolist]
command=java -jar /var/todolist-back-end/todo-list-api-1.0-SNAPSHOT.jar
autorstart=true
autorestart=true
stdout_logfile=/var/supervisor.log